# bin/mainsh
#!/bin/sh
### BEGIN INIT INFO
# Provides: mainsh
# Required-Start: $network
# Required-Stop: $network
# Default-Start: 2 3 5
# Default-Stop:
# Description: Starts main.sh
### END INIT INFO

case "$1" in
'start')
        sudo -u root /usr/bin/screen -S darknetstartup -d -m bash -c "/opt/darknet/main.sh start"
        ;;
'stop')
        sudo -u root /opt/darknet/main.sh "stop"
        ;;
*)
        echo "Usage: $0 { start | stop }"
        ;;
esac
exit 0
